# files/etc/supervisor/conf.d/cyanite.conf
[program:cyanite]
command = /opt/cyanite/cyanite.sh
logfile = /var/log/cyanite.log
autostart = true
autorestart = true 
