# file/etc/supervisor/conf.d/cyanite.conf
[program:cyanite]
command = /opt/cyanite/cyanite.sh
autostart = true
autorestart = true 
